# crontab
*/1 * * * * python3 /task.py >> /logfile
