# crontab
0/15 * * * * python3 /task.py >> /logfile
